# deploy/systemd/gatecontrol-gateway-update.service
[Unit]
Description=GateControl Gateway self-update (oneshot)

[Service]
Type=oneshot
ExecStart=/opt/gatecontrol-gateway/deploy/update.sh
