# deploy/systemd/gatecontrol-gateway-update.path
[Unit]
Description=Watch for GateControl Gateway pending-update flag

[Path]
PathExists=/state/pending-update
Unit=gatecontrol-gateway-update.service

[Install]
WantedBy=multi-user.target
